# docker/laravel-horizon/supervisord.d/laravel-horizon.conf.example
[program:laravel-horizon]
process_name=%(program_name)s_%(process_num)02d
user=sail
command=php /var/www/laravel/current/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
